# coredns/rootfs/etc/services.d/dnsdist/run
#!/usr/bin/with-contenv bashio
CONFIG="/etc/corefile.conf"

# Run coredns
bashio::log.info "Starting CoreDNS..."
exec coredns -conf "${CONFIG}" 